# roles/cinder-volume/templates/etc/init/cinder-volume.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start \
                       --chuid cinder \
                       --exec /usr/local/bin/cinder-volume \
                       -- --config-dir /etc/cinder
