# roles/network/templates/etc/init/neutron-lbaas-agent.conf
start on started neutron-openvswitch-agent
stop on stopping neutron-openvswitch-agent

respawn

exec start-stop-daemon --start \
                       --chuid neutron \
                       --make-pidfile --pidfile /var/run/neutron-lbaas-agent.pid \
                       --exec /usr/local/bin/neutron-lbaas-agent \
                       -- --config-file /etc/neutron/neutron.conf \
                          --config-file /etc/neutron/services/loadbalancer/haproxy/lbaas_agent.ini
