# roles/network/templates/etc/init/neutron-dhcp-agent.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start \
                       --chuid neutron \
                       --make-pidfile --pidfile /var/run/resource-agents/neutron-dhcp-agent.pid \ 
                       --exec /usr/local/bin/neutron-dhcp-agent \
                       -- --config-dir /etc/neutron \
                          --config-file /etc/neutron/neutron.conf \
                          --config-file /etc/neutron/dhcp_agent.ini
