# roles/network/templates/etc/init/neutron-dhcp-agent.conf
start on started neutron-openvswitch-agent
stop on stopping neutron-openvswitch-agent

respawn

exec start-stop-daemon --start \
                       --chuid neutron \
                       --exec /usr/local/bin/neutron-dhcp-agent \
                       -- --config-dir /etc/neutron \
                          --config-file /etc/neutron/neutron.conf \
                          --config-file /etc/neutron/dhcp_agent.ini
