# playbooks/neutron/templates/etc/init/neutron-server.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start \
                       --chuid neutron \
                       --exec /usr/local/bin/neutron-server \
                       -- --config-dir /etc/neutron \
                          --config-file /etc/neutron/neutron.conf \
                          --config-file /etc/neutron/plugins/openvswitch/ovs_neutron_plugin.ini
