# roles/logging/templates/etc/init/bbclog.conf
# bbclog - bbc logging daemon
#
# bbclog is a second rsyslog service
# used to forward logs.

description "bbc logging daemon"

start on filesystem
stop on runlevel [06]

expect fork
kill timeout 30
respawn

script
    exec /usr/local/sbin/bbclogd -f /etc/bbclog.conf -i /var/run/bbclog.pid
end script
