# roles/neutron-data-network/templates/etc/init/neutron-l3-agent.conf
{% if neutron.plugin == 'ml2' %}
start on started neutron-linuxbridge-agent
stop on stopping neutron-linuxbridge-agent
{% elif neutron.plugin == 'ovs' %}
start on started neutron-openvswitch-agent
stop on stopping neutron-openvswitch-agent
{% endif %}

respawn

exec start-stop-daemon --start \
                       --chuid neutron \
                       --make-pidfile --pidfile /var/run/neutron-l3-agent.pid \
                       --exec /usr/local/bin/neutron-l3-agent \
                       -- --config-dir /etc/neutron \
                          --config-file /etc/neutron/neutron.conf \
                          --config-file /etc/neutron/l3_agent.ini

