# playbooks/percona/templates/etc/init/garbd.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
expect fork

exec start-stop-daemon  --start \
                        --make-pidfile --pidfile /var/run/garbd.pid \
                        --exec /usr/bin/garbd -- -g mstack_db_cluster -l /var/log/garbd.log \
                        -a gcomm://
