# roles/common/templates/etc/profile.d/autologout.sh
TMOUT={{ common.idle_session_timeout }}
readonly TMOUT
export TMOUT
