# roles/controller/templates/etc/init/heat-api-cfn.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --chuid heat --exec /usr/local/bin/heat-api-cfn
