# roles/controller/templates/etc/init/neutron-server.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start \
                       --chuid neutron \
                       --exec /usr/local/bin/neutron-server \
                       -- --config-dir /etc/neutron \
                          --config-file /etc/neutron/neutron.conf \
{% if neutron.plugin == 'ml2' %}
                          --config-file /etc/neutron/plugins/ml2/ml2_plugin.ini \
{% elif neutron.plugin == 'ovs' %}
                          --config-file /etc/neutron/plugins/openvswitch/ovs_neutron_plugin.ini \
{% endif %}
{% if neutron.enable_lbaas %}
                          --config-file /etc/neutron/services/loadbalancer/haproxy/lbaas_agent.ini
{% endif %}
