# roles/keystone/templates/etc/init/keystone.conf
description "uwsgi for keystone"


start on runlevel [2345]
stop on runlevel [!2345]

respawn

pre-start script
  if [ ! -d /run/uwsgi ]; then
      mkdir /run/uwsgi/
      chown keystone /run/uwsgi
      chmod 775 /run/uwsgi
  fi
end script

post-stop script
  if [ -d /run/uwsgi ]; then
     rm -r /run/uwsgi
  fi
end script

exec {{ uwsgi_path }} --uid keystone --gid keystone --master --emperor /etc/keystone/uwsgi
